# rc.local
#!/bin/sh -e
#
# rc.local
#
# Print the IP address
_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi

if [ -x /bin/firewall.sh ]; then
	/bin/firewall.sh
fi

#/sbin/modprobe mbgclock
#/usr/sbin/mbgsvcd
/bin/ln -s /dev/ttyAMA0 /dev/refclock-1
#/bin/ln -s /dev/ttyUSB0 /dev/gps0
#/bin/sleep 90
#/usr/sbin/mbgsetsystime
/usr/sbin/service ntp start

exit 0


